# bckn-miner-continuous.service
[Unit]
Description=Bckn Continuous CPU Miner
After=network.target

[Service]
Type=simple
User=root
Environment="PRIVATE_KEY=YOUR_PRIVATE_KEY_HERE"
ExecStart=/bin/bash -c 'while true; do /usr/local/bin/bckn-miner-cpu ${PRIVATE_KEY}; sleep 3; done'
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Performance optimizations
CPUSchedulingPolicy=batch
Nice=-10

[Install]
WantedBy=multi-user.target